# deploy/vps/systemd/hallsync-backup.service
[Unit]
Description=Hall Sync verified database backup
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/var/www/hallsync/backend
EnvironmentFile=/etc/hallsync-api.env
ExecStart=/usr/bin/npm run production:backup:scheduled
User=www-data
Group=www-data
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/var/www/hallsync/backups

